# templates/etc/systemd/system/shairport-sync.target.j2
[Unit]
Description=shairport-sync Workers
Requires={% for sink in airplay_sinks %}shairport-sync_{{sink}}.service {% endfor %}

[Install]
WantedBy=multi-user.target